# nas/infra/systemd/carvis-nas.service
[Unit]
Description=Carvis NAS remote control server
After=network-online.target carvis-electron.service
Wants=network-online.target

[Service]
WorkingDirectory=%h/carvis-remote-smoke/nas
ExecStart=%h/carvis-remote-smoke/nas/carvis-nas-server
Restart=on-failure
RestartSec=2
Environment=CARVIS_NAS_CONFIG_DIR=%h/carvis-remote-smoke/nas/config
Environment=CARVIS_ELECTRON_API_URL=http://127.0.0.1:45932
Environment=CARVIS_OUTPUT_ROOT=%h/carvis-remote-smoke/output/runs
Environment=CARVIS_HISTORY_ROOT=%h/carvis-remote-smoke/output/runs
# Override this in a drop-in if the LAN IP or nginx domain changes.
Environment=CARVIS_NAS_PUBLIC_URL=http://127.0.0.1:8765

[Install]
WantedBy=default.target
